# Backend/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/15 * * * * cd /app && python etl.py >> /var/log/etl.log 2>&1
