# Backend/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
PYTHONPATH=/app
DATABASE_URL=postgresql://admin:admin@postgres:5432/admin

*/2 * * * * cd /app && python3 etl.py >> /var/log/etl.log 2>&1

# Empty line (required)

